# assets/etc/cron.d/fetchmail
# /etc/cron.d/fetchmail
0,15,30,45 * * * * /var/mail/fetchmail.pl > /dev/null
