# assets/etc/cron.d/fetchmail
# /etc/cron.d/fetchmail
*/15 * * * * /var/mail/fetchmail.pl > /dev/null
