# assets/etc/cron.d/fetchmail
# /etc/cron.d/fetchmail
*/5 * * * * root /var/mail/fetchmail.pl > /dev/null
